# setup/etc/systemd/system/kresd@.service.d/override.conf
[Unit]
Before=antizapret.service

[Service]
LimitNPROC=infinity